# shlls/remote
### BEGIN INIT INFO
# Provides:          chrome remote desktop
# Required-Start:    $time
# Required-Stop:     $time
# Default-Start:     5
# Default-Stop:      1 2 3 4 6
# Short-Description: Run chrome remote desktop
# Description: Start chrome remote desktop as user  
#              
#              
### END INIT INFO

#!/bin/zsh

status=$(HOME=/home/nr2 /opt/google/chrome-remote-desktop/chrome-remote-desktop --check-running)

case "$1" in
	start)
#	if [[ $status ]]; then
		echo starting chrome >> /tmp/remotetest
		su nr2 -c "HOME=/home/nr2 /opt/google/chrome-remote-desktop/chrome-remote-desktop --start"
#	fi
		;;
	stop)
#	if [[ ! $status ]]; then
		echo stopping chrome >> /tmp/remotetest
		su nr2 -c "HOME=/home/nr2 /opt/google/chrome-remote-desktop/chrome-remote-desktop --stop"
#	fi
		;;
	*)
		;;
esac
		
exit 0
